# supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=/usr/local/bin/python3 /usr/src/app/main.py

[program:mongodb]
command=/usr/bin/mongod
